# _scripts/coinrat-sync-bittrex-usd_btc.service
[Unit]
Description=Coinrat-Sync-bittrex-USD_BTC

[Service]
ExecStart=/home/coinrat/coinrat/coinrat-sync-bittrex-usd_btc.sh 
User=coinrat
